# supervisor/postgresql.conf
[supervisord]
nodaemon=true

[program:postgresql]
command=/opt/conda/bin/postgres -D /var/lib/postgresql
user=rdkit
autorestart=true
